# files/etc/profile.d/30-sysinfo.sh.bak
#!/bin/bash

export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

THIS_SCRIPT="sysinfo"
MOTD_DISABLE=""

machine_model=$(cat /proc/device-tree/model | tr -d "\000")
time=$(uptime | tr -d ',' | awk -F 'up ' '{print $NF}' | awk -F '  load' '{print $1}' | sed "s/  / /g")
if [[ -n $(echo -e "${time}" | grep "00") ]]; then
    time=$(echo "${time%:*} hour")
elif [[ -n $(echo -e "${time}" | grep ":") ]]; then
    time=$(echo -e "${time} min" | sed "s/:/ hour /")
fi
cpu_temp=$(awk '{printf ($0/100)}' /sys/class/thermal/thermal_zone0/temp)
load_avg=$(uptime | tr -d ',' | awk -F "average: " '{print $2}')
mem_used_total=$(free | grep "^Mem" | awk '{printf ("%d/%dMB",$3/1024,$2/1024)}')
boot_used_total=$(df -h /boot | awk '/\// {print $5"/"$2}')
root_used_total=$(df -h / | awk '/\// {print $5"/"$2}')
data_used_total=$(df -h /mnt/mmcblk2p4 | awk '/\// {print $5"/"$2}')

echo -e "Model: ${machine_model} , Uptime: ${time}"
echo -e "Load: ${load_avg}/${cpu_temp:0:2}.${cpu_temp:2:1}°C , Mem: ${mem_used_total}"
echo -e "Disk: Boot ${boot_used_total} Root ${root_used_total} Data ${data_used_total}"
echo ""
